# gunicorn.conf
[program:gunicorn]
directory=/home/ec2-user/flask-example
command=gunicorn main:app -b 0.0.0.0:5000 -c /home/ec2-user/flask-example/gunicorn.conf
user=ec2-user
autostart=true
autorestart=true
